# supervisord.conf
[supervisord]
nodaemon=true

[program:postgresql]
command=sudo /postgresql-entrypoint.sh

[program:tomcat6]
command=sudo /tomcat6-entrypoint.sh

[program:mapnik]
command=sudo /mapnik-entrypoint.sh


